# docker/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F --allow-to-run-as-root
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm_err.log

[program:telegram-poller]
command=php /hleb/console poll/telegram-task
directory=/hleb
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/telegram-poller.log
stderr_logfile=/var/log/supervisor/telegram-poller_err.log
user=www-data ; It's good practice to run commands as a non-root user if possible.
             ; Ensure this user has write permissions for any logs/files the script might create.
             ; The /hleb directory is mounted from the host, so permissions there matter.
             ; If issues arise, you might need to adjust this or ensure directory permissions.

; Note on user: The base php:fpm image runs php-fpm as root by default to allow binding to privileged ports,
; but then drops privileges to www-data for its worker processes.
; For the poller, using www-data is a good default if it has necessary file access.
; If your HLEB console command or the SDK needs to write files in areas www-data can't,
; you might need to adjust permissions or run as root (less ideal).
; The `php-fpm` command uses --allow-to-run-as-root because the main php-fpm process starts as root.
; Supervisor itself will also run as root in the container.
